# Sigma/06 Privilege Escalation/[T1053.005]Cron_Job_Privilege_Escalation.yml
title: Cron Job Creation for Privilege Escalation
technique: T1053.005
id: 4a5b6c7d-8e9f-0123-4567-89012345abcd
status: experimental
description: Detects suspicious cron job creation for privilege escalation
author: Infinit3i
date: 2024/04/27
references:
    - https://attack.mitre.org/techniques/T1053/005/
logsource:
    product: linux
    service: auditd
detection:
    selection_crontab:
        type: 'SYSCALL'
        comm: 'crontab'
        success: 'yes'
    selection_file_write:
        type: 'PATH'
        name|contains:
            - '/etc/crontab'
            - '/var/spool/cron'
            - '/etc/cron.d'
    selection_suspicious:
        comm: 'crontab'
        exe|contains:
            - '/tmp/'
            - '/var/tmp/'
            - '/dev/shm/'
    condition: (selection_crontab and selection_file_write) or
               selection_suspicious
falsepositives:
    - Legitimate system administration
    - Automated backup scripts
level: medium
tags:
  # MITRE ATT&CK tactic
  - attack.privilege_escalation
  # MITRE ATT&CK technique
  - attack.t1053.005
  # Platform/tool tags
  - windows
  - sysmon
  - attack.privilege_escalation
  - attack.t1053.005
  - attack.persistence
  - auditd
  - linux
