# Sigma/06 Privilege Escalation/[T1047]WMI_Event_Subscription_Abuse.yml
title: WMI Event Subscription Abuse for Privilege Escalation
technique: T1047
id: a1b2c3d4-e5f6-7890-abcd-ef1234567890
status: experimental
description: Detects WMI event subscription creation for privilege escalation and persistence
author: Infinit3i
date: 2026/04/27
references:
    - https://attack.mitre.org/techniques/T1047/
    - https://www.fireeye.com/blog/threat-research/2016/08/wmi_vs_wmi_monitor.html
logsource:
    product: windows
    service: sysmon
detection:
    selection_event_consumer:
        EventID: 1
        ParentImage|endswith: '\wmiprvse.exe'
        CommandLine|contains:
            - '__EventConsumer'
            - 'CommandLineEventConsumer'
            - 'ActiveScriptEventConsumer'
            - 'LogFileEventConsumer'
    selection_event_filter:
        EventID: 1
        Image|endswith: '\wmic.exe'
        CommandLine|contains|all:
            - '__EventFilter'
            - 'CREATE'
    selection_binding:
        EventID: 1
        CommandLine|contains|all:
            - '__FilterToConsumerBinding'
            - 'Consumer'
            - 'Filter'
    selection_privileged_execution:
        EventID: 1
        ParentImage|endswith: '\wmiprvse.exe'
        IntegrityLevel: 'System'
        CommandLine|contains:
            - 'whoami'
            - 'net user'
            - 'reg add'
            - 'schtasks'
    condition: (selection_event_consumer or selection_event_filter or
               selection_binding) or selection_privileged_execution
falsepositives:
    - Legitimate WMI monitoring solutions
    - System management tools
    - Security software with WMI components
level: high
tags:
  # MITRE ATT&CK tactic
  - attack.privilege_escalation
  # MITRE ATT&CK technique
  - attack.t1047
  # Platform/tool tags
  - windows
  - sysmon
  - attack.privilege_escalation
  - attack.t1047
  - attack.persistence
  - attack.t1546.003
  - sysmon
  - windows
