# Sigma/06 Privilege Escalation/[T1543.003]Scheduled_Task_Privilege_Escalation.yml
title: Scheduled Task Creation for Privilege Escalation
technique: T1543.003
id: 6c7d8e9f-0a1b-2345-6789-abcdef012345
status: experimental
description: Detects scheduled task creation with elevated privileges
author: Infinit3i
date: 2024/04/27
references:
  - https://attack.mitre.org/techniques/T1543/003/
logsource:
  product: windows
  service: sysmon
detection:
  selection_schtasks:
  EventID:
  Image|endswith:
  CommandLine|contains:
      - '/ru system'
      - '/rl highest'
      - 'NT AUTHORITY\SYSTEM'
  selection_taskschd:
  EventID:
  ParentImage|endswith:
  selection_registry:
  EventID:
  TargetObject|contains:
  condition: selection_schtasks or selection_taskschd or selection_registry
falsepositives:
  - Legitimate system administration
  - Software installation tasks
level: medium
tags:
  - attack.privilege_escalation
  - attack.t1543.003
  - windows
  - sysmon
  - attack.persistence
