# Sigma/05 Persistence/[T1053.005]ScheduledTaskXMLManipulation.yml
---
title: Scheduled Task XML Manipulation
technique: T1053.005
id: c3d4e5f6-7890-1234-5678-90abcdef1234
status: experimental
description: Detects direct manipulation of scheduled task XML files for persistence evasion
author: Infinit3i
date: 2026-04-27
references:
    - https://attack.mitre.org/techniques/T1053/005/
logsource:
    product: windows
    service: sysmon
detection:
    file_creation:
        EventID: 11
        TargetFilename|contains: '\Windows\System32\Tasks\'
        TargetFilename|endswith: '.xml'
    file_modification:
        EventID: 2
        TargetFilename|contains: '\Windows\System32\Tasks\'
    xml_tools:
        EventID: 1
        CommandLine|contains:
            - 'schtasks /xml'
            - 'Export-ScheduledTask'
            - 'Import-ScheduledTask'
    direct_xml_access:
        EventID: 1
        Image|endswith:
            - '\notepad.exe'
            - '\wordpad.exe'
            - '\powershell.exe'
        CommandLine|contains: '\Windows\System32\Tasks\'
    condition: (file_creation or file_modification) and
               (xml_tools or direct_xml_access)
fields:
    - TargetFilename
    - Image
    - CommandLine
    - ProcessId
    - User
    - ComputerName
falsepositives:
    - System maintenance tasks
    - Task Scheduler management tools
    - Backup and restore operations
level: medium
tags:
  # MITRE ATT&CK tactic
  - attack.persistence
  # MITRE ATT&CK technique
  - attack.t1053.005
  # Platform/tool tags
  - windows
  - sysmon
  - attack.persistence
  - attack.defense_evasion
  - attack.t1053
  - attack.t1053.005
  - windows
  - sysmon
