# Sigma/05 Persistence/[T1547.001]Registry_Run_Keys_Persistence.yml
title: Registry Run Keys Persistence Detection
id: a1b2c3d4-e5f6-7890-1234-567890abcdef
status: experimental
description: Detects suspicious process execution from Registry Run keys through registry modification events
author: Infinit3i
date: 2026/04/27
references:
    - https://attack.mitre.org/techniques/T1547/001/
logsource:
    product: windows
    service: sysmon
detection:
    selection_reg:
        EventID: 13
        TargetObject|contains:
            - '\Software\Microsoft\Windows\CurrentVersion\Run\'
            - '\Software\Microsoft\Windows\CurrentVersion\RunOnce\'
    selection_suspicious:
        - Details|contains:
            - 'powershell'
            - 'cmd.exe'
            - 'wscript'
            - 'cscript'
            - 'rundll32'
            - 'regsvr32'
        - Details|endswith:
            - '.exe -h'
            - '.exe --hidden'
            - '.bat'
            - '.vbs'
            - '.js'
            - '.ps1'
    condition: selection_reg and selection_suspicious
falsepositives:
    - Legitimate software installations
    - Administrative scripts
level: medium
tags:
  - attack.persistence
  - attack.t1547
  - attack.t1547.001
  - sysmon
  - windows
