# Sigma/05 Persistence/[T1547.001]BOF_Registry_Persistence.yml
---
title: BOF Registry-Based Persistence Mechanism
technique: T1547.001
id: d8e9f0a1-2b3c-4567-8901-bcdef2345678
status: experimental
description: Detects BOF registry persistence through autorun keys and COM hijacking
author: matthewiver@protonmail.com
date: 2026/04/27
references:
  - https://attack.mitre.org/techniques/T1547/001/
  - https://www.cobaltstrike.com/blog/beacon-object-files-luser-defined-reflective-loader/
logsource:
  product: windows
  service: sysmon
detection:
  autorun_keys:
    EventID: 1
    TargetObject|contains:
      - '\SOFTWARE\Microsoft\Windows\CurrentVersion\Run'
      - '\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce'
      - '\SOFTWARE\Microsoft\Windows\CurrentVersion\RunServices'
  com_hijack:
    EventID: 1
    TargetObject|contains:
      - '\SOFTWARE\Classes\CLSID'
      - '\SOFTWARE\Classes\Interface'
    Details|contains:
      - 'rundll32'
      - 'regsvr32'
  suspicious_values:
    EventID: 1
    Details|contains:
      - '.tmp'
      - '.dat'
      - 'powershell -e'
      - 'cmd /c'
  condition: (autorun_keys or com_hijack) and suspicious_values
fields:
  - TargetObject
  - Details
  - ProcessId
  - Image
  - EventType
falsepositives:
  - Software installation
  - System administration tools
  - Legitimate startup applications
level: high
tags:
  - attack.persistence
  - attack.t1547.001
  - bof
  - registry
  - autorun
  - com_hijacking
  - sysmon
  - windows
