# Sigma/05 Persistence/[T1547.001]Registry_Run_Keys_Execution_Monitor.yml
title: Registry Run Keys Execution Monitor
id: b2c3d4e5-f6g7-8901-2345-678901bcdefg
status: experimental
description: >-
  Detects process execution initiated from Registry Run keys through parent-child process relationships
author: Infinit3i
date: 2026/04/27
references:
  - https://attack.mitre.org/techniques/T1547/001/
logsource:
  product: windows
  service: sysmon
detection:
  selection_parent:
  EventID:
  ParentImage|endswith:
  selection_cmdline:
  CommandLine|contains:
      - 'HKLM\Software\Microsoft\Windows\CurrentVersion\Run'
      - 'HKCU\Software\Microsoft\Windows\CurrentVersion\Run'
      - 'CurrentVersion\Run'
  selection_suspicious:
      - Image|contains:
      - '\AppData\'
      - '\Temp\'
      - '\Users\Public\'
      - '\ProgramData\'
      - CommandLine|contains:
      - ' -enc '
      - ' -e '
      - 'powershell'
      - 'cmd /c'
      - 'rundll32'
  condition: selection_parent and (selection_cmdline or selection_suspicious)
falsepositives:
  - Legitimate startup applications
  - Software update mechanisms
level: high
tags:
  - attack.persistence
  - attack.t1547
  - attack.t1547.001
  - sysmon
  - windows
