# build/sentinel/Executable or Script Dropped in Persistence Path.json
{
  "kind": "Scheduled",
  "properties": {
    "displayName": "Executable or Script Dropped in Persistence Path",
    "description": "Detects creation of executables or script files in suspicious persistence paths like AppData or Startup folders using Sysmon Event ID 11.",
    "severity": "High",
    "enabled": false,
    "query": "Event\n| where TimeGenerated > ago(24h)\n| where (TargetFilename|contains has \"\\AppData\\\" or TargetFilename|contains has \"\\Start Menu\\Programs\\Startup\\\") and (TargetFilename|endswith has \".exe\" or TargetFilename|endswith has \".bat\" or TargetFilename|endswith has \".ps1\" or TargetFilename|endswith has \".vbs\")",
    "queryFrequency": "PT15M",
    "queryPeriod": "PT24H",
    "triggerOperator": "GreaterThan",
    "triggerThreshold": 0,
    "suppressionDuration": "PT1H",
    "suppressionEnabled": false,
    "alertRuleTemplateName": null,
    "incidentConfiguration": {
      "createIncident": true,
      "groupingConfiguration": {
        "enabled": false,
        "reopenClosedIncident": false,
        "lookbackDuration": "PT5H",
        "matchingMethod": "AllEntities",
        "groupByEntities": [],
        "groupByAlertDetails": [],
        "groupByCustomDetails": []
      }
    },
    "eventGroupingSettings": {
      "aggregationKind": "SingleAlert"
    },
    "alertDetailsOverride": null,
    "customDetails": {
      "SigmaRuleId": "9e212d9c-028d-4930-a186-19bdfde8ff87",
      "SigmaLevel": "high",
      "SigmaAuthor": "Matthew Iverson",
      "SigmaDate": "2025/04/20"
    },
    "entityMappings": [
      {
        "entityType": "Account",
        "fieldMappings": [
          {
            "identifier": "Name",
            "columnName": "Account"
          }
        ]
      },
      {
        "entityType": "Host",
        "fieldMappings": [
          {
            "identifier": "HostName",
            "columnName": "Computer"
          }
        ]
      }
    ],
    "sentinelEntitiesMappings": null,
    "templateVersion": null,
    "tactics": [],
    "techniques": [
      "T1547",
      "T001"
    ]
  }
}