# build/sentinel/Unauthorized Cron Job Creation.json
{
  "kind": "Scheduled",
  "properties": {
    "displayName": "Unauthorized Cron Job Creation",
    "description": "Detects syslog events indicating new or modified cron jobs by monitoring for messages such as 'crontab: installing new crontab' or 'crontab: installing for', which may signal unauthorized cron job creation.",
    "severity": "Medium",
    "enabled": false,
    "query": "Syslog\n| where TimeGenerated > ago(24h)\n| where (Message|contains has \"crontab: installing new crontab\" or Message|contains has \"crontab: installing for\")",
    "queryFrequency": "PT15M",
    "queryPeriod": "PT24H",
    "triggerOperator": "GreaterThan",
    "triggerThreshold": 0,
    "suppressionDuration": "PT1H",
    "suppressionEnabled": false,
    "alertRuleTemplateName": null,
    "incidentConfiguration": {
      "createIncident": true,
      "groupingConfiguration": {
        "enabled": false,
        "reopenClosedIncident": false,
        "lookbackDuration": "PT5H",
        "matchingMethod": "AllEntities",
        "groupByEntities": [],
        "groupByAlertDetails": [],
        "groupByCustomDetails": []
      }
    },
    "eventGroupingSettings": {
      "aggregationKind": "SingleAlert"
    },
    "alertDetailsOverride": null,
    "customDetails": {
      "SigmaRuleId": "e9c71683-79f8-40c0-93e4-1234567890ab",
      "SigmaLevel": "medium",
      "SigmaAuthor": "Matthew Iverson",
      "SigmaDate": "2025/04/08"
    },
    "entityMappings": [
      {
        "entityType": "Account",
        "fieldMappings": [
          {
            "identifier": "Name",
            "columnName": "Account"
          }
        ]
      },
      {
        "entityType": "Host",
        "fieldMappings": [
          {
            "identifier": "HostName",
            "columnName": "Computer"
          }
        ]
      }
    ],
    "sentinelEntitiesMappings": null,
    "templateVersion": null,
    "tactics": [],
    "techniques": [
      "T1053"
    ]
  }
}